# deploy/diana-bot.service.example
[Unit]
Description=Diana Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/diana-bot
ExecStart=/opt/diana-bot/.venv/bin/python main.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
